# ops/LaunchAgents/ai.bnimasta.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.bnimasta.backup</string>
    <key>Comment</key>
    <string>Daily backup of the BNI-Masta agent — vault, openclaw config, secrets, cloudflared, LaunchAgents. Runs 03:00 local.</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>3</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>~/.openclaw/agents/bni-masta/scripts/backup.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>~/.openclaw/agents/bni-masta/scripts/backup.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>~/.openclaw/agents/bni-masta/scripts/backup.stderr.log</string>
  </dict>
</plist>
